# fluxion/fluxion/rootfs/etc/s6-overlay/s6-rc.d/fluxion/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: FluxION ECS
# Runs the FluxION ECS main service
# ==============================================================================

bashio::log.info "Starting FluxION ECS service..."

exec /usr/bin/fluxion.sh
